# unturned-update-status/crontab
* * * * * root bash /root/update.sh >> /var/log/cron.log 2>&1
# Mandatory blank line